# contrib/remux.service
[Unit]
Description=ReMux Launcher

[Service]
ExecStart=/home/root/harmony/remux.exe
Restart=on-failure

[Install]
WantedBy=multi-user.target

